# bigAGI/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "Starting bigAGI ..."

## Run your program
exec /app/node_modules/.bin/next
